# conf/upstart.conf
description "Education Site"
start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 10 5

script
    /var/www/callowaysite/current_deploy/start_gunicorn.sh wsgi:application
end script
